# dot_config/systemd/user/hypridle.service
[Unit]
Description=hypridle idle/lock (SL1C3D-L4BS)
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/hypridle
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
