# dot_config/systemd/user/theme-propagation.service
[Unit]
Description=SL1C3D-L4BS theme propagation (generate theme.json at session start)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env sl1c3d theme apply

[Install]
WantedBy=default.target
