# dot_config/systemd/user/theme-propagation.service
[Unit]
Description=SL1C3D-L4BS theme propagation (generate theme.json at session start)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env sl1c3d theme apply
# Doctor severity: degraded (theme.json may be missing until next apply)

[Install]
WantedBy=graphical-session.target
