# dot_config/systemd/user/wallpaper-restore.service
[Unit]
Description=Waypaper restore last wallpaper (SL1C3D-L4BS)
After=graphical-session.target swww.service
PartOf=graphical-session.target
Wants=swww.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/waypaper --restore
# Doctor severity: degraded

[Install]
WantedBy=graphical-session.target
