# dot_config/SL1C3D-L4BS/system-config/getty@tty1.service.d/autologin.conf
# SL1C3D-L4BS — TTY1 autologin
# Install: see ../install-getty-autologin.sh

[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --autologin the_architect - ${TERM}
